# templates/etc/profile.d/tmout.sh.j2
{{ file_managed_by_ansible }}

#!/bin/bash

{% if rhel_09_412035 %}
declare -xr TMOUT=600
{% endif %}
